# deploy/imageopt-api.service
[Unit]
Description=Image Optimizer API (uvicorn)
After=network.target mysql.service redis.service

[Service]
User=www-data
Group=www-data
WorkingDirectory=/var/www/image-optimizer
EnvironmentFile=/var/www/image-optimizer/.env
ExecStart=/var/www/image-optimizer/venv/bin/uvicorn app.main:app \
    --host 127.0.0.1 \
    --port 8090 \
    --workers 4 \
    --log-level info
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
