# deploy/imageopt-beat.service
[Unit]
Description=Image Optimizer Celery Beat (scheduler)
After=network-online.target mysql.service redis-server.service
Wants=network-online.target mysql.service redis-server.service
StartLimitIntervalSec=0

[Service]
User=quailshack-image-optimizer
Group=quailshack-image-optimizer
WorkingDirectory=/home/quailshack-image-optimizer/htdocs/image-optimizer.quailshack.com
EnvironmentFile=/home/quailshack-image-optimizer/htdocs/image-optimizer.quailshack.com/.env
ExecStart=/home/quailshack-image-optimizer/htdocs/image-optimizer.quailshack.com/venv/bin/celery \
    -A app.workers.celery_app beat \
    --loglevel=info
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
